# offlineimap/crontab
*/30 * * * *    /usr/bin/offlineimap -q -u basic -c /email/offlineimaprc
11 */4 * * *    /usr/bin/offlineimap    -u basic -c /email/offlineimaprc
